# AMSL_ros_pkg/roomba_robot_meiji/roomba_500driver_meiji/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/amsl/AMSL_ros_pkg/roomba_robot_meiji/roomba_500driver_meiji/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
